# docker/supervisor/conf.d/50-scheduler.conf
[program:laravel-scheduler]
; Laravel scheduler runs every minute — shell loop avoids needing cron
command      = /bin/sh -c "while true; do php /var/www/html/artisan schedule:run --no-interaction >> /dev/null 2>&1; sleep 60; done"
directory    = /var/www/html
user         = www-data
priority     = 50
autostart    = true
autorestart  = true
stdout_logfile = /dev/null
stderr_logfile = /dev/null
